# inittab
# This is the "inittab", a configuration file used by the init
# supervisor to figure out what to run.
#
# Comments start with a '#' character and continue to the end of
# the line.

# Blank lines are ignored.

# Each non-blank, non-comment line must be the absolute path to an
# executable script or binary program to execute.  These scripts /
# binaries are expected to run in the foreground (i.e. no daemons)
# and are not expected to exit.
#
# If any of the child processes *do* exit, it is the job of init
# to restart them in a timely fashion.
#
# You cannot pass arguments to these scripts / binaries.  If you
# have a program that _needs_ arguments (nginx in the foreground
# must be run as `nginx -g "daemon off;"`, for instance), you will
# need to write a small wrapper script.
#
# Likewise with environment variables.
#
# These constraints simplify the code of the supervisor greatly.
#
# init does not handle logging at all.  All input/output streams
# (stdin, stdout, and stderr) are redirected to/from /dev/null,
# as appropriate.
#

/Users/jhunt/code/init/scripts/bash-server
/Users/jhunt/code/init/scripts/verse
